# systemd/openups.service
[Unit]
Description=OpenUPS Network Monitor with Auto-Shutdown
Documentation=https://github.com/surpasslight12/openups
After=network-online.target
Wants=network-online.target

# Service
Type=notify
# Optional delay before starting the main binary
ExecStartPre=/bin/sleep 300
ExecStart=/usr/local/bin/openups
TimeoutStartSec=360
Restart=on-failure
RestartSec=10
WatchdogSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openups
SyslogLevelPrefix=yes

# Environment
EnvironmentFile=-/etc/default/openups
Environment="OPENUPS_TARGET=1.1.1.1"
Environment="OPENUPS_INTERVAL=10"
Environment="OPENUPS_THRESHOLD=5"
Environment="OPENUPS_DRY_RUN=yes"
Environment="OPENUPS_LOG_LEVEL=info"
Environment="OPENUPS_TIMESTAMP=no"

# Security
CapabilityBoundingSet=CAP_NET_RAW CAP_SYS_BOOT
AmbientCapabilities=CAP_NET_RAW CAP_SYS_BOOT
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/openups

# 资源限制
MemoryMax=50M
TasksMax=10

[Install]
WantedBy=multi-user.target
